# crontab
* * * * * root "cd /satis && php bin/satis build satis.json /var/www/html"
