# cron/crontab
* * * * * root /satis/bin/satis build /schema/satis.json /usr/share/nginx/html

